# supervisord.conf
[supervisord]
nodaemon=true
user=root


[program:dapperdox]
command=sh /app/run.Unix.sh
